# deploy/sentinel-dashboard.service
[Unit]
Description=C-Sentinel Dashboard
Documentation=https://github.com/williamofai/c-sentinel
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
ExecStart=/opt/sentinel-dashboard/venv/bin/gunicorn -b 127.0.0.1:5000 app:app
WorkingDirectory=/opt/sentinel-dashboard
User=sentinel
Group=sentinel

# Environment - CHANGE THESE VALUES

# Database connection
Environment="DB_HOST=localhost"
Environment="DB_PORT=5432"
Environment="DB_NAME=sentinel"
Environment="DB_USER=sentinel"
Environment="DB_PASSWORD=CHANGE_ME_your_database_password"

# API key for sentinel agents to POST data
# Generate with: openssl rand -hex 16
Environment="SENTINEL_API_KEY=CHANGE_ME_generate_with_openssl_rand_hex_16"

# Flask session secret
# Generate with: python3 -c "import secrets; print(secrets.token_hex(32))"
Environment="FLASK_SECRET_KEY=CHANGE_ME_generate_with_secrets_token_hex_32"

# Dashboard login password hash
# Generate with: echo -n 'yourpassword' | sha256sum | cut -d' ' -f1
Environment="SENTINEL_ADMIN_PASSWORD_HASH=CHANGE_ME_sha256_hash_of_your_password"

# Dashboard URL (used in alert links)
Environment="DASHBOARD_URL=https://your-domain.com"

# Email alerting (optional - set ALERT_EMAIL_ENABLED=false to disable)
Environment="ALERT_EMAIL_ENABLED=false"
Environment="ALERT_SMTP_HOST=smtp.gmail.com"
Environment="ALERT_SMTP_PORT=587"
Environment="ALERT_SMTP_USER=your-email@gmail.com"
Environment="ALERT_SMTP_PASS=your-app-password"
Environment="ALERT_FROM=your-email@gmail.com"
Environment="ALERT_TO=alerts@your-domain.com"
Environment="ALERT_COOLDOWN_MINS=60"

# Slack alerting (optional - set ALERT_SLACK_ENABLED=false to disable)
# Get webhook URL from: https://api.slack.com/apps → Create App → Incoming Webhooks
Environment="ALERT_SLACK_ENABLED=false"
Environment="ALERT_SLACK_WEBHOOK=https://hooks.slack.com/services/YOUR/WEBHOOK/URL"
Environment="ALERT_SLACK_USERNAME=C-Sentinel"
# Optional: override default webhook channel
# Environment="ALERT_SLACK_CHANNEL=#security-alerts"

# Demo mode (optional - enables read-only public access without login)
# Useful for showcasing the dashboard publicly
Environment="DEMO_MODE=false"

# Restart policy
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectHome=read-only
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-dashboard

[Install]
WantedBy=multi-user.target

