# deploy/sentinel.service
[Unit]
Description=C-Sentinel System Monitor with Audit Integration
Documentation=https://github.com/williamofai/c-sentinel
After=network.target auditd.service
Wants=auditd.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sentinel --json --network --audit
User=root
Group=root

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Allow reading audit logs
ReadOnlyPaths=/var/log/audit

# Allow writing baseline data
ReadWritePaths=/var/lib/sentinel

# Exit codes 0, 1, 2 are all "success" (only 3 = error triggers restart)
SuccessExitStatus=0 1 2

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel

[Install]
WantedBy=multi-user.target
